# etc/xdg/autostart/Bojata GUI.desktop
[Desktop Entry]
Exec=lxterminal --command "/bin/bash -c 'sudo nice -n -20 /usr/bin/python3 /home/pi/bojata/bojata_gui.py'"
