# etc/xdg/autostart/bojata.desktop
[Desktop Entry]
Type=Application
Name=Bojata
Comment=Bojata autostart
Exec=lxterminal -e "/home/pi/bojata/scripts/run.sh 0"
